# root/var/lib/systemd/system/displaydaemon.service
[Unit]
Description=Display service

[Service]
ExecStart=python /opt/DisplayDaemon/main.py
Type=simple

[Install]
WantedBy=multi-user.target