# root/var/lib/systemd/system/displaydaemon.service
[Unit]
ExecStart=python /opt/DisplayDaemon/main.py
Type=simple