# root/var/lib/systemd/system/powerdaemon.service
[Unit]
ExecStart=python /opt/Power/powerdaemon.py
Type=simple