# root/etc/systemd/system/displaydaemon.service
[Unit]
Description=Display service

[Service]
ExecStart=/usr/bin/python /opt/DisplayDaemon/main.py
Type=simple

[Install]
WantedBy=multi-user.target